# update-motd.d/00-ponies
#! /bin/sh
echo "==============================================="

echo "Welcome $USER" | figlet | ponysay
